# k8s/packer/backup/packer/configs/docker-compose
#!/sbin/openrc-run

name="docker-bootstrap"


depend() {
    need docker
}

start_pre() {
    cd /docker
    docker-compose up -d
    
}